# config/storm-lv.conf
[program:storm-lv]
command=/opt/storm/apache-storm-0.10.0/bin/storm logviewer
stopsignal=INT
autorestart=true
autostart=true
startretries=3
startsecs=10
user=root

