# config/storm.conf
[program:storm]
command=/opt/storm/apache-storm-0.10.0/bin/storm supervisor
stopsignal=INT
autorestart=true
autostart=true
startretries=3
startsecs=10
user=root

