# config/storm-ui.conf
[program:storm-ui]
command=/opt/storm/apache-storm-0.10.0/bin/storm ui
stopsignal=INT
autorestart=true
autostart=true
startretries=3
startsecs=5
user=root

