# config/zookeeper.conf
[program:zookeeper]
command=/opt/zookeeper/zookeeper-3.4.7/bin/zkServer.sh start-foreground
stopsignal=INT
autorestart=true
autostart=true
startretries=3
startsecs=5
user=root

